# services/temporal/temporal-nginx.container
[Unit]
Description=Nginx Reverse Proxy for Temporal
Requires=temporal-ui.service
After=temporal-ui.service

[Container]
Image=docker.io/library/nginx:1.27-alpine
ContainerName=temporal-nginx
Network=temporal-network

# Volume mounts for configuration and certificates
Volume=systemd-nginx-config:/etc/nginx/conf.d:ro,Z
Volume=systemd-nginx-ssl:/etc/nginx/ssl:ro,Z

# Expose HTTP and HTTPS ports
PublishPort=80:80
PublishPort=443:443

# Health check
HealthCmd=curl -f http://temporal-server:7233/health || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target