# services/temporal/temporal-server.container
[Unit]
Description=Temporal Server
Requires=temporal-postgres.service
After=temporal-postgres.service

[Container]
Image=docker.io/temporalio/auto-setup:1.28.1
ContainerName=temporal-server
Network=temporal-network
PublishPort=7233:7233

# PostgreSQL configuration (connect to postgres user's database via localhost)
Environment=DB=postgres12
Environment=DB_PORT=5432
Environment=POSTGRES_USER=temporal
Environment=POSTGRES_PWD=temporal
Environment=POSTGRES_SEEDS=temporal-postgresql
Environment=POSTGRES_DB=temporal

# Basic server configuration - using host networking

# Volume mounts for dynamic configuration
Volume=temporal-config-volume:/etc/temporal/config/dynamicconfig

# Port 7233 automatically available with host networking

# Health check
HealthCmd=temporal --address localhost:7233 operator cluster health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target